# dht/dna-dht-bootstrap.service
[Unit]
Description=DNA Messenger DHT Bootstrap Node
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/dna-messenger/dht/build
Environment="LD_LIBRARY_PATH=/usr/local/lib"
ExecStart=/opt/dna-messenger/dht/build/persistent_bootstrap
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
